# linux/.config/systemd/user/bsu-login.timer
[Unit]
Description=Run BSu Login periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=12h
Persistent=true

[Install]
WantedBy=timers.target
